# proxmox/coredns-zone-sync.timer
[Unit]
Description=Periodically sync svc-150/153 venue IPs into the CoreDNS scion.zone

[Timer]
# ~3 min after boot (give CT215/216/217 time to start + DHCP/SLAAC), then
# every 5 min. Each run is a cheap no-op unless a venue address (or the
# CoreDNS deployment) changed (idempotent script), so it catches a venue
# move / lease change / late CoreDNS deploy hands-free. Persistent so a
# missed run fires on the next boot.
OnBootSec=3min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
