# proxmox/coredns-zone-sync.service
[Unit]
# Keeps the `web`/`web2` A/AAAA records in CT216's (svc-152) CoreDNS zone in
# sync with the CURRENT venue address of svc-150 (CT215) and svc-153
# (CT217). Their eth1 leg is DHCP/SLAAC, so the address changes on a
# rebuild/lease renewal/venue move; this re-derives it and rewrites the zone
# ONLY when it actually changed (the script is idempotent). CoreDNS's `file`
# plugin polls and reloads the zone on its own — no restart triggered here.
Description=Sync svc-150/153 venue IPs into the CoreDNS scion.zone (CT216)
After=pve-guests.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
Environment=AUTO=1
ExecStart=/usr/bin/bash /root/ietf-scion-testbed/proxmox/coredns-zone-sync.sh
