# linkd/pkg/prerm
#!/bin/sh
set -e
systemctl disable --now scion-linkd.service || true
