# linkd/pkg/scion-linkd.service
[Unit]
Description=SCION testbed link shaping daemon
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/usr/bin/scion-linkd -config /etc/scion-linkd/config.toml
Restart=on-failure
RestartSec=2
CapabilityBoundingSet=CAP_NET_ADMIN
AmbientCapabilities=CAP_NET_ADMIN
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/etc/scion
ProtectHome=yes

[Install]
WantedBy=multi-user.target
