# config/supervisord.conf
[program:br1-150-1]
autostart = false
autorestart = false
environment = TZ=UTC,GODEBUG="cgocheck=0"
stdout_logfile = logs/br1-150-1.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/router --config gen/AS150/br1-150-1.toml

[program:cs1-150-1]
autostart = false
autorestart = false
environment = TZ=UTC
stdout_logfile = logs/cs1-150-1.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/control --config gen/AS150/cs1-150-1.toml

[program:sd1-150]
autostart = false
autorestart = false
environment = TZ=UTC
stdout_logfile = logs/sd1-150.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/daemon --config gen/AS150/sd.toml

[group:as1-150]
programs = br1-150-1,cs1-150-1,sd1-150

[program:br1-151-1]
autostart = false
autorestart = false
environment = TZ=UTC,GODEBUG="cgocheck=0"
stdout_logfile = logs/br1-151-1.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/router --config gen/AS151/br1-151-1.toml

[program:cs1-151-1]
autostart = false
autorestart = false
environment = TZ=UTC
stdout_logfile = logs/cs1-151-1.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/control --config gen/AS151/cs1-151-1.toml

[program:sd1-151]
autostart = false
autorestart = false
environment = TZ=UTC
stdout_logfile = logs/sd1-151.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/daemon --config gen/AS151/sd.toml

[group:as1-151]
programs = br1-151-1,cs1-151-1,sd1-151

[program:br1-152-1]
autostart = false
autorestart = false
environment = TZ=UTC,GODEBUG="cgocheck=0"
stdout_logfile = logs/br1-152-1.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/router --config gen/AS152/br1-152-1.toml

[program:cs1-152-1]
autostart = false
autorestart = false
environment = TZ=UTC
stdout_logfile = logs/cs1-152-1.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/control --config gen/AS152/cs1-152-1.toml

[program:sd1-152]
autostart = false
autorestart = false
environment = TZ=UTC
stdout_logfile = logs/sd1-152.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/daemon --config gen/AS152/sd.toml

[group:as1-152]
programs = br1-152-1,cs1-152-1,sd1-152

[program:br1-153-1]
autostart = false
autorestart = false
environment = TZ=UTC,GODEBUG="cgocheck=0"
stdout_logfile = logs/br1-153-1.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/router --config gen/AS153/br1-153-1.toml

[program:cs1-153-1]
autostart = false
autorestart = false
environment = TZ=UTC
stdout_logfile = logs/cs1-153-1.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/control --config gen/AS153/cs1-153-1.toml

[program:sd1-153]
autostart = false
autorestart = false
environment = TZ=UTC
stdout_logfile = logs/sd1-153.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/daemon --config gen/AS153/sd.toml

[group:as1-153]
programs = br1-153-1,cs1-153-1,sd1-153

[program:br1-154-1]
autostart = false
autorestart = false
environment = TZ=UTC,GODEBUG="cgocheck=0"
stdout_logfile = logs/br1-154-1.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/router --config gen/AS154/br1-154-1.toml

[program:cs1-154-1]
autostart = false
autorestart = false
environment = TZ=UTC
stdout_logfile = logs/cs1-154-1.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/control --config gen/AS154/cs1-154-1.toml

[program:sd1-154]
autostart = false
autorestart = false
environment = TZ=UTC
stdout_logfile = logs/sd1-154.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/daemon --config gen/AS154/sd.toml

[group:as1-154]
programs = br1-154-1,cs1-154-1,sd1-154

[program:br1-155-1]
autostart = false
autorestart = false
environment = TZ=UTC,GODEBUG="cgocheck=0"
stdout_logfile = logs/br1-155-1.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/router --config gen/AS155/br1-155-1.toml

[program:cs1-155-1]
autostart = false
autorestart = false
environment = TZ=UTC
stdout_logfile = logs/cs1-155-1.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/control --config gen/AS155/cs1-155-1.toml

[program:sd1-155]
autostart = false
autorestart = false
environment = TZ=UTC
stdout_logfile = logs/sd1-155.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/daemon --config gen/AS155/sd.toml

[group:as1-155]
programs = br1-155-1,cs1-155-1,sd1-155

[program:br1-156-1]
autostart = false
autorestart = false
environment = TZ=UTC,GODEBUG="cgocheck=0"
stdout_logfile = logs/br1-156-1.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/router --config gen/AS156/br1-156-1.toml

[program:cs1-156-1]
autostart = false
autorestart = false
environment = TZ=UTC
stdout_logfile = logs/cs1-156-1.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/control --config gen/AS156/cs1-156-1.toml

[program:sd1-156]
autostart = false
autorestart = false
environment = TZ=UTC
stdout_logfile = logs/sd1-156.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/daemon --config gen/AS156/sd.toml

[group:as1-156]
programs = br1-156-1,cs1-156-1,sd1-156

[program:br1-157-1]
autostart = false
autorestart = false
environment = TZ=UTC,GODEBUG="cgocheck=0"
stdout_logfile = logs/br1-157-1.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/router --config gen/AS157/br1-157-1.toml

[program:cs1-157-1]
autostart = false
autorestart = false
environment = TZ=UTC
stdout_logfile = logs/cs1-157-1.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/control --config gen/AS157/cs1-157-1.toml

[program:sd1-157]
autostart = false
autorestart = false
environment = TZ=UTC
stdout_logfile = logs/sd1-157.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/daemon --config gen/AS157/sd.toml

[group:as1-157]
programs = br1-157-1,cs1-157-1,sd1-157

[program:br1-158-1]
autostart = false
autorestart = false
environment = TZ=UTC,GODEBUG="cgocheck=0"
stdout_logfile = logs/br1-158-1.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/router --config gen/AS158/br1-158-1.toml

[program:cs1-158-1]
autostart = false
autorestart = false
environment = TZ=UTC
stdout_logfile = logs/cs1-158-1.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/control --config gen/AS158/cs1-158-1.toml

[program:sd1-158]
autostart = false
autorestart = false
environment = TZ=UTC
stdout_logfile = logs/sd1-158.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/daemon --config gen/AS158/sd.toml

[group:as1-158]
programs = br1-158-1,cs1-158-1,sd1-158

[program:br1-159-1]
autostart = false
autorestart = false
environment = TZ=UTC,GODEBUG="cgocheck=0"
stdout_logfile = logs/br1-159-1.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/router --config gen/AS159/br1-159-1.toml

[program:cs1-159-1]
autostart = false
autorestart = false
environment = TZ=UTC
stdout_logfile = logs/cs1-159-1.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/control --config gen/AS159/cs1-159-1.toml

[program:sd1-159]
autostart = false
autorestart = false
environment = TZ=UTC
stdout_logfile = logs/sd1-159.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/daemon --config gen/AS159/sd.toml

[group:as1-159]
programs = br1-159-1,cs1-159-1,sd1-159

[program:br1-160-1]
autostart = false
autorestart = false
environment = TZ=UTC,GODEBUG="cgocheck=0"
stdout_logfile = logs/br1-160-1.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/router --config gen/AS160/br1-160-1.toml

[program:cs1-160-1]
autostart = false
autorestart = false
environment = TZ=UTC
stdout_logfile = logs/cs1-160-1.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/control --config gen/AS160/cs1-160-1.toml

[program:sd1-160]
autostart = false
autorestart = false
environment = TZ=UTC
stdout_logfile = logs/sd1-160.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/daemon --config gen/AS160/sd.toml

[group:as1-160]
programs = br1-160-1,cs1-160-1,sd1-160

[program:br1-161-1]
autostart = false
autorestart = false
environment = TZ=UTC,GODEBUG="cgocheck=0"
stdout_logfile = logs/br1-161-1.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/router --config gen/AS161/br1-161-1.toml

[program:cs1-161-1]
autostart = false
autorestart = false
environment = TZ=UTC
stdout_logfile = logs/cs1-161-1.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/control --config gen/AS161/cs1-161-1.toml

[program:sd1-161]
autostart = false
autorestart = false
environment = TZ=UTC
stdout_logfile = logs/sd1-161.log
redirect_stderr = True
startretries = 0
startsecs = 5
priority = 100
command = bin/daemon --config gen/AS161/sd.toml

[group:as1-161]
programs = br1-161-1,cs1-161-1,sd1-161

[program:dispatcher]
autostart = false
autorestart = false
environment = TZ=UTC
stdout_logfile = logs/dispatcher.log
redirect_stderr = True
startretries = 0
startsecs = 1
priority = 50
command = bin/dispatcher --config gen/dispatcher/disp.toml

