# ansible/templates/coredns/coredns.service.j2
[Unit]
Description=CoreDNS (scion. TLD resolver + Quad9 forwarder)
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
ExecStart=/usr/local/bin/coredns -conf /etc/coredns/Corefile
DynamicUser=yes
AmbientCapabilities=CAP_NET_BIND_SERVICE
ReadOnlyPaths=/etc/coredns
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
