# ansible/templates/playground/ttyd.service.j2
[Unit]
Description=ttyd web terminal (SCION playground)
After=network-online.target scion-sciond.service
Wants=network-online.target

[Service]
# %I is not used; ExecStart is templated by ansible with credential + bind.
ExecStart=/usr/local/bin/ttyd --port 7681 --interface {{ pubnet_addr }} \
  --writable --credential scion:{{ booth_code }} --max-clients 30 \
  --client-option titleFixed="SCION Playground AS1-{{ playground_as }}" \
  /usr/local/bin/guest-login.sh
Restart=on-failure
User=root

[Install]
WantedBy=multi-user.target
